# docker/php/crontab
# Lance la gestion des tâches planifiées dans Symfony toutes les minutes.
*	*	*	*	*	cd /app && php bin/console shapecode:cron:run
